# 09-MagicMirror/config_files/supervisor/mirror_webapp.conf
[program:mirror_webapp]
environment=MAGIC_MIRROR_SETTINGS="/home/pi/mirror_app.conf"
command = gunicorn -b 127.0.0.1:5000 mirror:app
directory = /home/pi/mirror_webapp
user = pi
