# 07-SecuritySystem/config_files/supervisor/security-webapp.conf
[program:security_webapp]
environment=SECURITY_APP_SETTINGS="/home/pi/security.conf"
command = gunicorn -b 127.0.0.1:5000 security:app
directory = /home/pi/security_webapp
user = pi
