# 08-RemoteReboticArm/config_files/supervisor/robot_arm_webapp.conf
[program:robot_arm_webapp]
environment=ROBOT_ARM_SETTINGS="/home/pi/robot_arm.conf"
command = gunicorn -b 127.0.0.1:5000 robot_arm:app
directory = /home/pi/robot_arm_webapp
user = pi
